# external/supervisor/mairy-backend.conf
; mairy-backend queue workers (local) — map prod mairy.conf → LocalStack FIFO queues
; Prod: /var/www/39mail + TM-39mail-Prod-SQS-*.fifo
; Local: /var/www/mairy-backend + queue_*.fifo (khớp .env SQS_* và LocalStack init)
; create/send worker cùng queue_default vì app dùng 1 env SQS_CREATE_MESSAGES / SQS_SEND_MESSAGES

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:default-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs-fifo --queue=queue_default.fifo --sleep=3 --timeout=180 --tries=2
directory=/var/www/mairy-backend
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/default-worker.log
stopwaitsecs=3600

[program:create-messages-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs-fifo --queue=queue_default.fifo --sleep=3 --timeout=180 --tries=1
directory=/var/www/mairy-backend
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/create-messages-worker.log
stopwaitsecs=3600

[program:send-messages-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs-fifo --queue=queue_default.fifo --sleep=3 --timeout=180
directory=/var/www/mairy-backend
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/send-messages-worker.log
stopwaitsecs=3600

[program:send-log-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs-fifo --queue=queue_send_log.fifo --sleep=3 --timeout=180 --tries=2
directory=/var/www/mairy-backend
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/send-log-worker.log
stopwaitsecs=3600

[program:crawl-task-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs-fifo --queue=queue_crawl_task.fifo --sleep=3 --timeout=180 --tries=2
directory=/var/www/mairy-backend
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crawl-task-worker.log
stopwaitsecs=3600

[program:crawl-task-nat-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs-fifo --queue=queue_crawl_nat_task.fifo --sleep=3 --timeout=180 --tries=2
directory=/var/www/mairy-backend
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crawl-task-nat-worker.log
stopwaitsecs=3600
